# iptables-watcher.service
[Unit]
Description=iptables watcher service

[Service]
ExecStart=/usr/local/bin/iptables-watcher.sh

[Install]
WantedBy=multi-user.target
